# discord-blue.service
[Unit]
Description=Discord Blue
After=multi-user.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=discord-blue
Group=discord-blue
WorkingDirectory=/opt/discord-blue
Environment=HOME=/var/lib/discord-blue
ExecStart=/opt/discord-blue/.venv/bin/discord-blue
Restart=on-failure
RestartSec=15
TimeoutStopSec=10
StateDirectory=discord-blue
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
